# wslg-init.service
[Unit]
Description=Check and recreat wslg symlink

[Service]
ExecStart=/usr/bin/wslg-init.sh

[Install]
WantedBy=multi-user.target